# debian/postinst
#!/bin/sh
# postinst script for hcron
#
# see: dh_installdeb(1)

# Source debconf library.
. /usr/share/debconf/confmodule

set -e


# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
	# update-rc.d hcron defaults 89 11 >/dev/null

    # /var/spool/hcron
	if [ ! -d /var/spool/hcron ]; then
		mkdir -p /var/spool/hcron
	fi
	chmod a+rwx,o+t /var/spool/hcron

	# /var/log/hcron
	if [ ! -d /var/log/hcron ]; then
		mkdir -p /var/log/hcron
	fi
	chmod a+rx /var/log/hcron

	# /var/lib/hcron
	if [ ! -d /var/lib/hcron ]; then
		mkdir -p /var/lib/hcron
	fi
	chmod a+rx /var/lib/hcron

	if [ ! -d /var/lib/hcron/event_lists ]; then
		mkdir -p /var/lib/hcron/event_lists
	fi
	chmod a+rx /var/lib/hcron/event_lists

	if [ ! -d /var/lib/hcron/trees ]; then
		mkdir -p /var/lib/hcron/trees
	fi
	chmod go-rwx /var/lib/hcron/trees
    ;;


    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0